# rpi-image/layer/getty@tty1.service.d/autologin.conf.tpl
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin AUTOLOGIN_USER --noclear %I 38400 linux
